# system/refresher.service
[Unit]
Before=graphical.target

[Service]
Type=simple
ExecStart=/bin/refresher

[Install]
WantedBy=multi-user.target

