# rules.d/50-firmware.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/usr/local/sbin/firmware.sh"
